# roles/richard/templates/richard.uwsgi.conf.j2
description "Richard uWSGI"
start on level [2345]
stop on level [06]

setuid richard
setgid richard

exec uwsgi --ini {{ richard_dir }}/uwsgi.ini
start on startup
